# recipes-houseos/houseos-services/files/mnt-data.mount
[Unit]
Description=Data mount service unit

[Mount]
Type=ext4
What=/dev/sda5
Where=/mnt/data

[Install]
WantedBy=multi-user.target